# services/dashboard.service
[Unit]
Description=Telegram Monitoring Dashboard
After=network.target
Wants=network.target

[Service]
Type=simple
User=%USER%
WorkingDirectory=%DIR%/telegram-dashboard
ExecStartPre=/bin/bash -c 'test -d %DIR%/telegram-dashboard/.next || (cd %DIR%/telegram-dashboard && npm install && npm run build)'
ExecStart=%NPM% start
Restart=on-failure
RestartSec=10
StandardOutput=append:%DIR%/logs/dashboard.log
StandardError=append:%DIR%/logs/dashboard.log

[Install]
WantedBy=multi-user.target
